# .config/systemd/user/meowrch-pywal.service
[Unit]
Description=Apply pywal colors from the current meowrch wallpaper

[Service]
Type=oneshot
ExecStart=%h/bin/pywal-reload.sh
